# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
user=root

[program:postgres]
# Wait for initdb (run by nextjs-starter on first boot) before starting postgres.
# Use -k to set socket dir to the directory we created and chowned in Dockerfile.
command=sh -c "until [ -f /var/lib/postgresql/data/PG_VERSION ]; do sleep 1; done; exec su postgres -c 'postgres -D /var/lib/postgresql/data -k /run/postgresql'"
autorestart=true
priority=10
# Stream to Docker logs so failures are visible via `docker logs`
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextjs-starter]
command=/docker/healthcheck.sh
priority=20
autorestart=false
startsecs=0
# Stream to Docker logs so healthcheck.sh output is visible via `docker logs`
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
